# ansible-prom-graf/playbooks/roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=root
ExecStart={{ node_exporter_exec_start }}

[Install]
WantedBy=default.target
